# deploy/systemd/sigrok.service
# Sigrok Discord bot — systemd (system unit)
#
# Install:
#   sudo install -m 644 deploy/systemd/sigrok.service /etc/systemd/system/sigrok.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now sigrok.service
#
# Before enable:
#   - Set User= / Group= to the Linux account that owns the repo (or create a dedicated one).
#   - Put the clone at ~<that-user>/Sigrok (see WorkingDirectory / ExecStart), or edit both
#     lines to match your layout (e.g. /opt/sigrok/... instead of %h/Sigrok/...).
#   - In-project venv (recommended): `poetry config virtualenvs.in-project true && poetry install`
#     so `.venv/bin/python` exists. Otherwise point ExecStart at `$(poetry env info -p)/bin/python`.

[Unit]
Description=Sigrok Discord bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sigrok
Group=sigrok
# %h expands to the home directory of User= (see systemd.unit(5) specifiers).
WorkingDirectory=%h/Sigrok

# Loads settings.toml + .secrets.toml from WorkingDirectory
ExecStart=%h/Sigrok/.venv/bin/python %h/Sigrok/src/sigrok/bot.py

Restart=on-failure
RestartSec=10

# Hardening (optional; remove lines that break your setup)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
